# overlay-variant-stumpy-moblab/chromeos-base/chromeos-bsp-moblab/files/init/moblab-devserver-cleanup-init.conf
# Copyright (c) 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Clean up the devserver static directory daily."
author        "chromium-os-dev@chromium.org"

start on started moblab-devserver-init

script
  if [ -z "${STATIC_DIR}" ] ; then
    logger -t "${UPSTART_JOB}" "No static directory param passed."
    exit 1
  fi
  while true; do
    logger -t "${UPSTART_JOB}" "Running devserver clean up."
    sudo -u moblab \
      /usr/local/autotest/site_utils/admin/clean_staged_images.py \
      -a 24 -p 24 "${STATIC_DIR}" || true
    sleep 86400
  done
end script
