# overlay-variant-stumpy-moblab/chromeos-base/chromeos-bsp-moblab/files/init/moblab-gsoffloader_s-init.conf
# Copyright (c) 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Kick off the gs_offloader to clean up the hosts folder."
author        "chromium-os-dev@chromium.org"

start on stopped autotest_setup
respawn

exec sudo -u moblab /usr/local/autotest/site_utils/gs_offloader.py -s -o --days_old=14
