# project-moblab/chromeos-base/chromeos-bsp-moblab/files/init/moblab-containers-dir-init.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Set up the containers directory on the external storage."
author        "chromium-os-dev@chromium.org"

start on stopped moblab-external-storage-init RESULT=ok

env CONTAINER_DIR=/mnt/moblab/containers
env AUTOTEST_CONTAINER_DIR=/usr/local/autotest/containers

script
  mkdir -p "${CONTAINER_DIR}"
  ln -s "${CONTAINER_DIR}" "${AUTOTEST_CONTAINER_DIR}"
end script