# overlay-x86-mario/chromeos-base/chromeos-bsp-x86-mario/files/99-light-sensor.rules
# Calibrate the light sensor when the tsl2563 driver is installed.
ACTION=="add", SUBSYSTEM=="drivers", KERNEL=="tsl2563", RUN+="light-sensor-set-multiplier.sh"
