# overlay-variant-stumpy-moblab/chromeos-base/chromeos-bsp-moblab/files/moblab-scheduler-init.conf
# Copyright (c) 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start the autotest scheduler."
author        "chromium-os-dev@chromium.org"

start on stopped moblab-apache-init
respawn

env HOME=/root
exec /usr/local/autotest/scheduler/monitor_db.py /usr/local/autotest/results
