# project-moblab/chromeos-base/chromeos-bsp-moblab/files/moblab-base-container-init.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Setup server side packaging base container."
author        "chromium-os-dev@chromium.org"

start on stopped moblab-containers-dir-init

script
  # Sleep a second after moblab-containers-dir-init has stopped.
  sleep 1
  logger -t "${UPSTART_JOB}" "Setting up base container."
  python /usr/local/autotest/site_utils/lxc.py -s
  chmod -R +r /usr/local/autotest/containers/base
end script
